# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=service cron restart
autorestart=false